# systemd/robust-config-watch.service
[Unit]
Description=Robust Config Sync Daemon — watches config file and broadcasts D-Bus deltas
Documentation=https://github.com/http418imateapot/robust-config-exchange
After=dbus.service
Requires=dbus.service

[Service]
Type=notify
# Run as a dedicated non-privileged user in production.
# Change this to a dedicated service account (e.g. robust-config) as needed.
User=root
ExecStart=/usr/local/bin/robust_config \
    --config /etc/robust-config/config.conf \
    --bus system \
    --log-level info \
    watch
Restart=on-failure
RestartSec=5s

# Systemd watchdog: kill and restart if no WATCHDOG=1 ping within 30 s
WatchdogSec=30s
NotifyAccess=main

# Hardening
ProtectSystem=strict
ReadWritePaths=/etc/robust-config
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
